# portbridge/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# Launch nginx in the foreground (s6 supervises this process)
exec nginx -g 'daemon off;'
